# roles/monitoring/templates/svscan.service.j2
[Unit]
Description=Svscan

[Service]
ExecStart=/usr/bin/svscan /services

[Install]
WantedBy=multi-user.target
